# systemd/hopeturtle-gps.timer
[Unit]
Description=Run HopeTurtle GPS logger every 5 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
