# systemd/hopeturtle-gps.service
#      ___________    _____
#  /               \ |  o  |
# |                |/   __\|
# |  _______________  /
#   |_|_|     |_|_|
#
# HopeTurtle GPS Service
# Runs gps_snapshot.py once (systemd timer calls this every 5 minutes).
# Reads GPS from soft-serial GPIO17 by default and logs CSV + JSON prototype.
# After each run, also displays abbreviated fix info on OLED.

[Unit]
Description=HopeTurtle GPS one-shot logger

[Service]
Type=oneshot

# ---- Environment configuration ----
Environment=HT_GPS_MODE=soft
Environment=HT_GPS_SOFT_RX=17
Environment=HT_BAUD=9600
Environment=HT_DATA_DIR=/home/hopeturtle002/hopeturtle/data

# ---- Working directory ----
WorkingDirectory=/home/hopeturtle002/hopeturtle

# ---- Main command ----
ExecStart=/usr/bin/python3 /home/hopeturtle002/hopeturtle/src/gps_snapshot.py

# ---- OLED feedback ----
ExecStartPost=/usr/bin/python3 /home/hopeturtle002/hopeturtle/src/oled_status.py brief

# ---- Run as user ----
User=hopeturtle002
Group=hopeturtle002

# ---- Lower priority ----
Nice=10
